# service.takealug.epg-grabber/upgrade-epg.cron
#!/bin/bash

20 5 * * * /local/data/src/service.takealug.epg-grabber/upgrade-epg.sh
